# etc/udev/rules.d/99-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHOROTY}, RUN+="~/Documents/scripts/monitor-hotplug.sh"
